# http_proxy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting Nginx web server..."
exec nginx -g "daemon off;"