# http_proxy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: HTTP Proxy
# Runs Nginx web server
# ==============================================================================
bashio::log.info "Starting Nginx web server..."
exec nginx -g "daemon off;"