# http_proxy/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: HTTP Proxy
# Runs API server
# ==============================================================================
bashio::log.info "Starting API server..."
cd /opt/api
exec node server.js