# http_proxy/rootfs/etc/services.d/api/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting API server..."
cd /opt/api
exec node server.js