# http_proxy/rootfs/etc/services.d/squid/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: HTTP Proxy
# Runs Squid proxy server
# ==============================================================================
bashio::log.info "Starting Squid proxy server..."
exec squid -N