# http_proxy/rootfs/etc/services.d/squid/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Add-on: HTTP Proxy
# Runs Tinyproxy proxy server
# ==============================================================================
bashio::log.info "Starting Tinyproxy proxy server..."
exec tinyproxy -d -c /etc/tinyproxy/tinyproxy.conf