# http_proxy/rootfs/etc/services.d/squid/run
#!/usr/bin/with-contenv bashio
set -e

bashio::log.info "Starting Squid proxy server..."
exec squid -N